# pi-monitor.service
[Unit]
Description=Pi Monitor - Continuous System Monitoring Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=abhinav
Group=abhinav
WorkingDirectory=/home/abhinav/pi-mon
ExecStart=/usr/bin/python3 /home/abhinav/pi-mon/backend/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-monitor

# Environment variables
Environment=PI_MONITOR_API_KEY=pi-monitor-api-key-2024
Environment=PYTHONUNBUFFERED=1

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/abhinav/pi-mon/backend /home/abhinav/pi-mon

# Auto-restart on failure
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
